# profile.d/F18_knowhowERP.sh
export KNOWHOW_ERP_ROOT=/opt/knowhowERP
export HARBOUR_ROOT=/opt/knowhowERP/hbout
PATH=$KNOWHOW_ERP_ROOT/bin:$KNOWHOW_ERP_ROOT/util:$KNOWHOW_ERP_ROOT/lib:$HARBOUR_ROOT/bin:$PATH
